# conf/supervisord.conf
[supervisord]
nodaemon=true
user=root
pidfile=/tmp/supervisord.pid

[program:app]
; ────────────────────────────────────────────────────────────────
; CRITICAL: Athena creds use ATHENA_ prefix (NOT AWS_ACCESS_KEY_ID)
;
; Why? The HealthService uses SSMClient which picks up the default
; AWS credential chain (AWS_ACCESS_KEY_ID env var). If we set
; AWS_ACCESS_KEY_ID to our Athena-only IAM user, SSM calls fail
; with AccessDeniedException and the app never reports HEALTHY.
;
; By using ATHENA_ACCESS_KEY_ID / ATHENA_SECRET_KEY:
;   - HealthService/SSMClient uses the TDP container IAM role → SSM works
;   - AthenaClient explicitly reads ATHENA_ACCESS_KEY_ID → Athena works
; ────────────────────────────────────────────────────────────────
command=/bin/sh -c "exec env NODE_ENV=production ENABLE_HEALTHCHECK=true ATHENA_ACCESS_KEY_ID=${AWS_ACCESS_KEY_ID} ATHENA_SECRET_KEY=${AWS_SECRET_ACCESS_KEY} ATHENA_REGION=us-east-1 ATHENA_WORKGROUP=takeda-gxp-dev ATHENA_DATABASE=takeda_gxp_dev__tss__external ATHENA_S3_OUTPUT_LOCATION=s3://takeda-dev-prod-athena-results/takeda_gxp_dev/ node packages/server/dist/index.js"
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopwaitsecs=30
stopsignal=TERM
killasgroup=true
stopasgroup=true
user=reactapp
